# deploy/chanlun.service
# ── BTC Chanlun Analyzer — systemd Service ──
# Deploy to: /etc/systemd/system/chanlun.service
#
# Usage:
#   sudo systemctl enable chanlun
#   sudo systemctl start chanlun
#   sudo systemctl status chanlun
#   journalctl -u chanlun -f

[Unit]
Description=BTC Chanlun Analyzer — FastAPI Backend
Documentation=https://github.com/ls569333469/BTC-PM
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/BTC-PM/backend
EnvironmentFile=/opt/BTC-PM/backend/.env

ExecStart=/opt/BTC-PM/backend/.venv/bin/python -m uvicorn \
    app.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 2 \
    --log-level info

# Auto-restart on failure
Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chanlun

[Install]
WantedBy=multi-user.target
